# launchd/com.user.brewdiff.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<!--
=============================================================================
com.user.brewdiff.plist
Spike & Chilli Home Network (E2) — Weekly Brew Diff Email Agent
=============================================================================

WHAT:   Fires brew-diff-email.sh weekly on Sunday at 04:00.
        Diffs Brewfile.NLMacMiniM1 vs Brewfile.NLMacbookProM3 from
        network-ops/data/brew-sync/ and emails a styled HTML summary.
        Sends a heartbeat "all in sync" email if no differences found.

NOT:    Does NOT install packages. Does NOT modify Brewfiles.

SCHEDULE: Sunday 04:00 — 2 hours after brew-sync (02:00) so Brewfiles
          are freshly dumped before the diff runs.

DEPLOY on NLMacMiniM1 (primary — always home):
  1. sudo cp scripts/brew-diff-email.sh /usr/local/bin/brew-diff-email.sh
     sudo chmod 755 /usr/local/bin/brew-diff-email.sh

  2. cp launchd/com.user.brewdiff.plist ~/Library/LaunchAgents/
     launchctl load ~/Library/LaunchAgents/com.user.brewdiff.plist

  3. Verify:
     launchctl list | grep brewdiff

MANUAL TEST (don't wait for Sunday 04:00):
  launchctl kickstart gui/$(id -u)/com.user.brewdiff
  tail -f /Volumes/network-ops/logs/brew_NLMacMiniM1_diff_*.log

RELOAD AFTER PLIST EDITS:
  launchctl unload ~/Library/LaunchAgents/com.user.brewdiff.plist
  launchctl load   ~/Library/LaunchAgents/com.user.brewdiff.plist

NOTE: First run requires approving Mail.app Automation permission:
  System Settings → Privacy & Security → Automation
=============================================================================
-->

<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.user.brewdiff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/bin/brew-diff-email.sh</string>
    </array>

    <!-- Fire weekly: Sunday (Weekday=0) at 04:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Prevent double-fire within 10 minutes -->
    <key>ThrottleInterval</key>
    <integer>600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/brewdiff.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/brewdiff.err</string>

</dict>
</plist>
